# config/systemd/ai-sandbox-mounts.service
[Unit]
Description=ai-sandbox host passthrough mounts
Documentation=file:/home/ai/ai-sandbox/config/ensure-sandbox-mounts.sh
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
RemainAfterExit=yes
EnvironmentFile=-/etc/ai-sandbox/target-user.env
ExecStart=/bin/sh -c '/home/${AI_SANDBOX_TARGET_USER:-ai}/ai-sandbox/config/ensure-sandbox-mounts.sh "${AI_SANDBOX_TARGET_USER:-ai}"'
ExecStop=/bin/sh -c 'umount /mnt/host-config /mnt/host-secrets /mnt/host-workspace 2>/dev/null; fusermount -u /mnt/host-ai-sandbox 2>/dev/null; umount /mnt/host-ai-sandbox 2>/dev/null; true'
Restart=on-failure
RestartSec=15s
TimeoutStartSec=120s

[Install]
WantedBy=multi-user.target
